# roles/packages/templates/downloadupdates/downloadupdates.timer
[Unit]
Description=Download Pacman Updates

[Timer]
OnCalendar=*-*-* 07:00:00
# Do not execute job if it missed a run due to machine being off (only run if system is powered ON)
Persistent=false
AccuracySec=10min

[Install]
WantedBy=timers.target